# rpi_flow/systemd/px4flow-rpi.service
[Unit]
Description=Raspberry Pi Optical Flow to PX4 (OPTICAL_FLOW_RAD)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/px4flow-rpi/rpi_flow
ExecStart=/opt/px4flow-rpi/rpi_flow/px4flow_rpi.sh /opt/px4flow-rpi/rpi_flow/config.json
Restart=on-failure
RestartSec=1
User=root

[Install]
WantedBy=multi-user.target
